# playbooks/roles/restic_schedule/templates/systemd/restic@.service.j2
[Unit]
Description=Restic backup for %I

[Service]
Type=oneshot
ExecStart={{ restic_schedule_install_path }}/restic-backup %i

[Install]
WantedBy=multi-user.target
